# pwn/proc-self-wtf-1/crontab
* * * * * /usr/bin/docker run --rm --network none --volume /private/flags/flag.txt:/flag.txt -d sleep
